# arch-config/.config/systemd/user/restic-backup-storage.timer
[Unit]
Description=Backup home with restic daily
[Timer]
OnCalendar=daily
Persistent=true
[Install]
WantedBy=timers.target
