# arch-config/.config/systemd/user/restic-backup-home.timer
[Unit]
Description=Backup home with restic weekly
[Timer]
OnCalendar=daily
Persistent=true
[Install]
WantedBy=timers.target
