# roles/config/files/etc/systemd/system/vmware-networks.service.d/override.conf
[Service]
ExecStartPost=/bin/sh -c "/usr/bin/chgrp vmware /dev/vmnet*"
ExecStartPost=/bin/sh -c "/usr/bin/chmod g+rw /dev/vmnet*"
